# configs/rclone-nextcloud.service
[Unit]
Description=Rclone Mount Nextcloud (WebDAV)
After=network.target
Wants=network.target

[Service]
Type=simple

# Garante diretórios e tenta desmontar mount antigo preso (sem falhar se não existir)
ExecStartPre=/usr/bin/mkdir -p %h/nextcloud %h/.cache/rclone-nextcloud
ExecStartPre=/bin/sh -lc 'mountpoint -q "%h/nextcloud" && fusermount3 -uz "%h/nextcloud" || true'

# Mount
# Foco: reduzir travadas quando a internet cai:
# - timeouts baixos (timeout/contimeout)
# - poucas tentativas (retries/low-level-retries)
# - VFS cache full para servir mais do cache
ExecStart=/usr/bin/rclone mount nextcloud: %h/nextcloud \
  --vfs-cache-mode full \
  --cache-dir %h/.cache/rclone-nextcloud \
  --vfs-cache-max-size 100G \
  --vfs-cache-max-age 24h \
  --vfs-write-back 10s \
  --vfs-read-ahead 128M \
  --dir-cache-time 30s \
  --attr-timeout 1s \
  --timeout 10s \
  --contimeout 5s \
  --low-level-retries 1 \
  --retries 1 \
  --retries-sleep 2s \
  --transfers 4 \
  --checkers 8 \
  --umask 022 \
  --allow-other \
  --log-level INFO \
  --log-file %h/.cache/rclone-nextcloud/rclone.log

# Remote Control (RC) - deixe comentado, caso queira ativar no futuro
# 1) Descomente as 2 linhas abaixo
# 2) Reinicie o serviço
# 3) Use: rclone rc core/version
#
#  --rc \
#  --rc-addr 127.0.0.1:5572

ExecStop=/bin/fusermount3 -uz %h/nextcloud

Restart=on-failure
RestartSec=10
TimeoutStopSec=20

[Install]
WantedBy=default.target

# nano ~/.config/systemd/user/rclone-nextcloud.service
# mkdir -p ~/nextcloud ~/.config/systemd/user ~/.cache/rclone-nextcloud
# systemctl --user daemon-reload && systemctl --user enable --now rclone-nextcloud.service

# TSHOOT

# systemctl --user daemon-reload
# systemctl --user restart rclone-nextcloud.service
# systemctl --user status rclone-nextcloud.service --no-pager -l

